# playbooks/templates/deploy-agent.service.j2
[Unit]
Description=Deploy Agent - Webhook service for container management
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=basecard
WorkingDirectory={{ user_home }}/src/{{ app_name }}
EnvironmentFile={{ user_home }}/src/{{ app_name }}/.env
ExecStart={{ user_home }}/src/{{ app_name }}/{{ app_binary | default('deploy-agent') }}
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deploy-agent

[Install]
WantedBy=multi-user.target
